# linux/custom-aur-packages/yazi-fm1-shim/yazi-fm1.service
[Unit]
Description=FileManager1 shim for yazi (opens Kitty + yazi)

[Service]
ExecStart=/usr/bin/yazi-fm1
Restart=on-failure

[Install]
WantedBy=default.target
